# battery_log.timer
[Unit]
Description=Logging battery charge

[Timer]
Unit=battery_log.service
OnBootSec=10
OnUnitInactiveSec=5
